# deploy/systemd/ai-story-backup.timer.example
[Unit]
Description=Run AI Story Adventure backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=ai-story-backup.service

[Install]
WantedBy=timers.target
